# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:varnishd]
command=/usr/sbin/varnishd -f /etc/varnish/default.vcl -a 0.0.0.0:$(ENV_HTTP_PORT)02d -s malloc,$(ENV_MAX_CACHE_SIZE)s -F
