# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:varnishd]
command=/bin/bash -c "/usr/sbin/varnishd -f ${VARNISH_CFG_DIR}/default.vcl -a 0.0.0.0:${HTTP_PORT} -s malloc,${MAX_CACHE_SIZE} -F"
